# src/GaIA/pkgs/atop/atop-1.16/atop.cron
# start atop daily at midnight
1 0 * * * root /var/opt/UNItools//etc/atop/atop.daily
