# src/GaIA/pkgs/atop/atop-1.19/atop.cron.dist
# start atop daily at midnight
1 0 * * * root /etc/atop/atop.daily
